# os-config/etc/systemd/system/first-boot.service
[Unit]
Description=One time boot script
[Service]
Type=simple
ExecStart=/usr/sbin/first-boot.sh
[Install]
WantedBy=multi-user.target
