# os-config/etc/systemd/user/flatpak-user-update.timer
[Unit]
Description=Update User Flatpaks Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target